# lib/itamae/plugin/recipe/god/templates/etc/systemd/system/god.service.erb
[Unit]
Description=God: A process monitoring framework in Ruby
After=network.target

[Service]
EnvironmentFile=-/etc/sysconfig/god
PIDFile=<%= @pid %>
ExecStart=<%= @god_bin %> -c <%= @config %> -P <%= @pid %> -l <%= @log %> --log-level <%= @log_level %>
ExecStop=<%= @god_bin %> quit

[Install]
WantedBy=multi-user.target
