# lib/itamae/plugin/recipe/god/templates/etc/init.d/god.erb
#!/bin/bash
#
# god       Startup script for God monitoring tool.
#
# chkconfig: - 85 15
# description: god monitors your system
#              

. /etc/rc.d/init.d/functions

PID=<%= @pid %>
CONFIG=<%= @config %>
LOG=<%= @log %>
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
RETVAL=0

case "$1" in
    start)
      echo -n 'Starting god: '
      god -P $PID -c $CONFIG -l $LOG --log-level <%= @log_level %>
      echo
      RETVAL=$?
  ;;
    stop)
      echo -n 'Stopping god: '
      god quit
      echo
      RETVAL=$?
  ;;
    restart)
      echo -n 'Stopping god: '
      god quit
      echo
      echo -n 'Starting god: '
      god -P $PID -c $CONFIG -l $LOG --log-level <%= @log_level %>
      echo
      RETVAL=$?
  ;;
    status)
      god status
      RETVAL=$?
  ;;
    *)
      echo "Usage: god {start|stop|restart|status}"
      exit 1
  ;;
esac

exit $RETVAL
